# meta-terminal-system-core/recipes-terminal-system/terminal-system-files/files/usb-automount@.service
[Unit]
Description=USB auto mount service

[Service]
Type=simple
Environment="SCRIPT_ARGS=%I"
ExecStart=/usr/bin/usb-automount.sh $SCRIPT_ARGS
