# meta-terminal-system-edgeplant/recipes-terminal-system/terminal-system-files/files/jasmine/ssd-mount@.service
[Unit]
Description=SSD mount service

[Service]
Type=oneshot
RemainAfterExit=yes
Environment="SCRIPT_ARGS=%I"
ExecStart=/usr/bin/ssd-mount.sh $SCRIPT_ARGS
